# blockchain-notary-supervisor.conf
[program:notary]
command=/usr/local/bin/gunicorn notary.wsgi:application --workers=3
directory=/var/www/blockchain-notary/
user=nobody
autostart=true
autorestart=true

